# etc/cron.d/etherest.conf
MAILTO=""
* * * * * root /sbin/service etherest watchdog
